# dot_config/systemd/user/my-upgrade-if-needed.timer
[Unit]
Description=Check if my-upgrade is needed (every day)

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true

[Install]
WantedBy=timers.target
